# roles/portal/files/portal_commands_receivable.socket
[Socket]
ListenStream=/var/run/commands_receivable.sock
RemoveOnStop=true

[Install]
WantedBy=sockets.target
